# dropbox-personal.service
[Unit]
Description=Start personal Dropbox instance after login
After=graphical-session.target network-online.target
Wants=network-online.target
PartOf=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h
Environment=HOME=%h
ExecStart=%h/scripts/sync-personal-dropbox.sh

[Install]
WantedBy=default.target
