# dropbox-work.service
[Unit]
Description=Mount work Dropbox remote after login
After=graphical-session.target xdg-desktop-autostart.target network-online.target
Wants=network-online.target xdg-desktop-autostart.target
PartOf=graphical-session.target

[Service]
WorkingDirectory=%h
Environment=HOME=%h
ExecStart=%h/scripts/rclone-mount.sh "Work Dropbox:" "Work Dropbox"

[Install]
WantedBy=graphical-session.target
