# dropbox-personal.service
[Unit]
Description=Mount personal Dropbox remote after login
After=graphical-session.target xdg-desktop-autostart.target network-online.target
Wants=network-online.target xdg-desktop-autostart.target
PartOf=graphical-session.target
StartLimitIntervalSec=1min
StartLimitBurst=3
Requires=dropbox-work.service
After=dropbox-work.service

[Service]
WorkingDirectory=%h
Environment=HOME=%h
ExecStart=%h/scripts/rclone-mount.sh "Personal Dropbox:" "Personal Dropbox"
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=graphical-session.target
